# launchd/com.claudecode.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claudecode.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/USER/.bun/bin/bun</string>
        <string>/Users/USER/.claude/daemon/main.ts</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/USER/.claude/daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/USER/.claude/daemon.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/USER</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/USER/.bun/bin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
